# systemd/system/man-db.service
[Unit]
Description=Update man-db

[Service]
ExecStart=/etc/cron.daily/man-db
Type=oneshot
