# systemd/system/man-db.timer
[Unit]
Description=Update man-db
Documentation=man:mandb(8)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
