# systemd/system/portage-sync.service
[Unit]
Description=daily tree sync

[Service]
ExecStart=emaint sync -a
Type=oneshot
